# buildroot/external/package/seeed-voicecard/files/S90seeed-voicecard
#!/bin/sh
#
# /etc/init.d/seeed-voicecard: Lanza el script de configuraci贸n de Seeed Voicecard al arranque.
#

### BEGIN INIT INFO
# Provides:          seeed-voicecard
# Required-Start:    $local_fs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Initialize Seeed Voicecard hardware.
### END INIT INFO

DAEMON="/usr/bin/seeed-voicecard"

case "$1" in
  start)
    echo "Starting Seeed Voicecard configuration..."
    # Ejecuta el script principal que realiza la detecci贸n y configuraci贸n.
    if [ -x "$DAEMON" ]; then
        $DAEMON
    fi
    ;;
  stop)
    # Este servicio es de tipo "one-shot", no hay proceso que detener.
    echo "Stop operation is not applicable for seeed-voicecard."
    ;;
  restart|force-reload)
    # Para reiniciar, simplemente se vuelve a ejecutar la l贸gica de inicio.
    $0 start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}" >&2
    exit 1
    ;;
esac

exit 0
